# demoapi.supervisor.conf
[program:demoapi]
command=/webapps/demoapi/demoapi/gunicorn.start
user=demoapi
autostart=true
autorestart=true